# DEBIAN/prerm
#!/bin/sh
set -ex

echo "Prerm fired"
update-rc.d -f domicube remove
