# DEBIAN/prerm
#!/bin/sh
set -ex

echo "Stopping domicube daemon"
/etc/init.d/domicube stop || true
echo "Removing InitScripts"
update-rc.d -f domicube remove
